# supervisor.x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -usepw -shared -rfbport 5900 -noncache
environment=DISPLAY=:1
priority=15
startsecs=3
startretries=4
